# app/crontab/root
*/1 * * * * /app/crontab/update_articles.py > app/crontab/log.txt
