# app/crontab/root
*/1 * * * * /app/crontab/update_articles.py
